# udev/70-logitech-trueforce.rules
# Logitech direct-drive racing wheels: hidraw + wheel_* sysfs access
#
# Grants the currently-logged-in session user (via systemd-logind
# uaccess) and anyone in the "input" group access to:
#   - /dev/hidraw*  for libtrueforce / userspace tooling
#   - /sys/class/hidraw/hidrawN/device/wheel_*  for Oversteer, the
#     in-tree sysfs knobs, and user scripts
#
# Without this rule the sysfs knobs are mode 0664 root:root and every
# wheel_* write needs sudo.
#
# Supported wheels:
#   046D:C262  Logitech G920 Racing Wheel
#   046D:C266  Logitech G923 (PlayStation/PC) - TrueForce-capable
#   046D:C26D  Logitech G923 (Xbox/PC, alt PID) - TrueForce-capable
#   046D:C26E  Logitech G923 (Xbox/PC) - TrueForce-capable
#   046D:C268  Logitech G PRO Racing Wheel (PlayStation/PC)
#   046D:C272  Logitech G PRO Racing Wheel (Xbox/PC)
#   046D:C276  Logitech RS50 Direct Drive (PS/PC)
#
# The G923 PIDs get hidraw access so Logitech's SDK DLLs (TrueForce)
# can reach the wheel under Proton the same way they do for the
# RS50/G PRO; the G923 speaks the same TrueForce stream protocol.

ACTION!="add|change", GOTO="rs50_end"
SUBSYSTEM!="hidraw", GOTO="rs50_end"

# One alternation instead of a GOTO per PID: adding a wheel is a
# five-character diff that cannot forget the GOTO.
ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c262|c266|c26d|c26e|c268|c272|c276", GOTO="rs50_apply"
GOTO="rs50_end"

LABEL="rs50_apply"

# /dev/hidraw* for libtrueforce and direct userspace tooling.
TAG+="uaccess", GROUP="input", MODE="0660"

# wheel_* sysfs attributes (only exist on the HID++ interface's hidraw
# node; the RUN invocation is a no-op on interfaces that lack them).
RUN+="/bin/sh -c 'for f in /sys/class/hidraw/%k/device/wheel_* /sys/class/hidraw/%k/device/range /sys/class/hidraw/%k/device/gain /sys/class/hidraw/%k/device/autocenter /sys/class/hidraw/%k/device/combine_pedals /sys/class/hidraw/%k/device/damper_level /sys/class/hidraw/%k/device/spring_level /sys/class/hidraw/%k/device/friction_level; do [ -e \"$f\" ] || continue; chgrp input \"$f\"; chmod g+w \"$f\"; done'"

LABEL="rs50_end"
